# dwb.service
[Unit]
Description=lightweight web browser dwb
After=wm.target

[Service]
ExecStart=/usr/bin/dwb

[Install]
WantedBy=mystuff.target
