# eufloria.service
[Unit]
Description=hib rts game
After=wm.target

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/eufloria

[Install]
WantedBy=wm.target
