# gimp.service
[Unit]
Description=Image editor
After=wm.target

[Service]
ExecStart=/usr/bin/gimp

[Install]
WantedBy=mystuff.target
