# luakit.service
[Unit]
Description=lightweight web browser luakit
After=wm.target

[Service]
ExecStart=/usr/bin/luakit

[Install]
WantedBy=mystuff.target
