# mbsync@.timer
[Unit]
Description=sync mailboxes with %i

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=mystuff.target
