# mcomix.service
[Unit]
Description=comic book viewer
After=wm.target

[Service]
ExecStart=/usr/bin/mcomix

[Install]
WantedBy=mystuff.target
