# minecraft.service
[Unit]
Description=Creative mining game
After=wm.target

[Service]
ExecStart=/usr/bin/minecraft
