# mpd.service
[Unit]
Description=Music Player Daemon
Before=tmux@.service

[Service]
ExecStart=/usr/bin/mpd --no-daemon

[Install]
WantedBy=mystuff.target
