# repose@.service
[Unit]
Description=Update the %i repo

[Service]
Type=simple
ExecStart=/usr/bin/repose -Ufc %h/.cache/pacman/pkg/%i.db
