# steam.service
[Unit]
Description=Steam launcher
After=wm.target

[Service]
ExecStart=/usr/bin/steam
ExecStop=/usr/bin/pkill -F %h/.steampid
ExecStop=/usr/bin/pkill -9 MainThrd

[Install]
WantedBy=wm.target
