# udiskie.service
[Unit]
Description=handle automounting

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=mystuff.target
