# urxvtd-user.service
[Unit]
Description=Urxvt Daemon

[Service]
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=mystuff.target
