# urxvtd-user.service
[Unit]
Description=Urxvt Daemon

[Service]
Environment=SHELL=/bin/zsh
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=mystuff.target
