# zathura.service
[Unit]
Description=pdf reader
After=wm.target

[Service]
ExecStart=/usr/bin/zathura

[Install]
WantedBy=mystuff.target
